# config/etc/profile.d/hadoop.sh
HADOOP_PREFIX=/vagrant/hadoop
export HADOOP_PREFIX
